# pitempd.service
[Unit]
Description=Ambient temperature daemon

[Service]
Type=notify
ExecStart=/home/pi/pitempd/pitempd

[Install]
WantedBy=multi-user.target
